# analysis-project/src/main/resources/docs/runtime-image/python-exec.conf
; P1-A 容器侧常驻 Python 执行服务进程守护
; 落到容器内路径: /etc/supervisor/conf.d/python-exec.conf
; Python 脚本:   /opt/runner/server.py
;
; 重启策略说明:
;   autorestart=true  → 意外退出自动拉起(crash / 超时强杀 pool 时不会退出)
;   startretries=3   → 连续 3 次失败后放弃(如 Python 版本不对)
;   stdout_logfile 走 JSON 行格式,与 JVM 日志关联 key=request_id

[program:python-exec]
command = python3 /opt/runner/server.py
directory = /opt/runner
user = root

; 自动重启 + 无需等待正常退出
autorestart = true
startretries = 3
startsecs = 5

; 环境变量注入 — 不改代码即可调优
environment =
    PYTHON_EXEC_HOST="127.0.0.1",
    PYTHON_EXEC_PORT="8000",
    PYTHON_EXEC_WORKERS="4",
    PYTHON_EXEC_DEFAULT_TIMEOUT="60",
    PYTHON_EXEC_MAX_TIMEOUT="300",
    PYTHON_EXEC_MAX_OUTPUT_BYTES="64000"

; 日志 — JOSN 行,可以和 JVM 日志一起用 request_id 关联
stdout_logfile = /var/log/python-exec/stdout.log
stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 5
stderr_logfile = /var/log/python-exec/stderr.log
stderr_logfile_maxbytes = 100MB
stderr_logfile_backups = 3

; 一旦 stderr 有输出就往上走(不要缓冲)
redirect_stderr = false